# srvbak.cron.sample
#!/bin/sh

set -e
export MAILER_SUBJECT="backup of ` hostname ` @ ` date +'%F %T' `"
/opt/src/srvbak/srvbak.bash /opt/src/srvbak/srvbakrc 2>&1 \
  | /opt/src/mailer/mailer.bash

# --
